# .config/autostart/keyswap.desktop
[Desktop Entry]
Name=KeySwap
Exec=xmodmap -e "keycode 66 = Escape" &&
Terminal=false
Type=Application
